# deploy/systemd/waterwall-proxy-restart.service
[Unit]
Description=Restart waterwall-proxy

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart waterwall-proxy.service
